# api-server.service
[Unit]
Description=API server
After=syslog.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/aps/service-prediction-api
Environment="APP_CONFIG=staging.env"
ExecStartPre=source activate aps
ExecStart=uwsgi uwsgi.ini
Restart=always

[Install]
WantedBy=multi-user.target
