# celery-worker.service
[Unit]
Description=Celery worker
After=syslog.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/aps/service-prediction-api
Environment="PYTHONOPTIMIZE=1"
Environment="APP_CONFIG=staging.env"
ExecStart=/opt/anaconda/envs/aps/bin/celery -A celery_worker.celery worker -E --loglevel=debug --concurrency=1 --max-tasks-per-child=1
Restart=always

[Install]
WantedBy=multi-user.target
