# uwsgi.conf
start on [2345]
stop on [06]

script
    exec uwsgi --ini ./uwsgi.ini
end script
